# files/etc/supervisor/conf.d/amavis.conf
[program:amavisd-new]
command=/usr/sbin/amavisd-new foreground
directory=/tmp
priority=500

