# files/etc/supervisor/conf.d/apache2.conf
[program:apache2]
command=/etc/supervisor/scripts/generic.sh '/usr/sbin/apache2ctl -k start' '/usr/sbin/apache2ctl stop'
directory=/tmp

