# files/etc/supervisor/conf.d/clamav-daemon.conf
[program:clamav-daemon]
command=/usr/sbin/clamd -c /etc/clamav/clamd.conf
startsecs=5
directory=/tmp
priority=400

