# files/etc/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
startsecs=5
directory=/tmp
priority=300

