# files/etc/supervisor/conf.d/postfix.conf
[program:postfix] 
command=/etc/supervisor/scripts/generic.sh '/usr/sbin/postfix -c /etc/postfix start' '/usr/sbin/postfix stop'
directory=/tmp
priority=300

