# files/etc/supervisor/conf.d/rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
directory=/tmp
priority=100

