# files/etc/supervisor/conf.d/zarafa-spooler.conf
[program:zarafa-spooler]
command=/usr/bin/zarafa-spooler -F
directory=/tmp

