# files/etc/supervisor/conf.d/zarafa-monitor.conf
[program:zarafa-monitor]
command=/usr/bin/zarafa-monitor -F
directory=/tmp

