# files/etc/supervisor/conf.d/zarafa-search.conf
[program:zarafa-search]
command=/usr/bin/zarafa-search -F
directory=/tmp

