# files/etc/supervisor/conf.d/zarafa-gateway.conf
[program:zarafa-gateway]
command=/usr/bin/zarafa-gateway -F
directory=/tmp

