# files/etc/supervisor/conf.d/zarafa-dagent.conf
[program:zarafa-dagent]
command=/usr/bin/zarafa-dagent -l
directory=/tmp

