# files/etc/supervisor/conf.d/zarafa-ical.conf
[program:zarafa-ical]
command=/usr/bin/zarafa-ical -F
directory=/tmp

